# src/ext2_SDpart_skel/etc/inittab
::sysinit:/etc/rc
::restart:/sbin/init
::shutdown:/etc/rc.shutdown
::reboot:/etc/rc.reboot
::ctrlaltdel:/etc/rc.reboot
